# tellsticklive/rootfs/etc/services.d/learn/finish
#!/command/with-contenv bashio
# ==============================================================================
# Learn service finish script - do not restart this one-shot service
# ==============================================================================
exec true
